# systemd/lauschomat-web.service
[Unit]
Description=Lauschomat Web Visualization Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lauschomat
Group=lauschomat
ExecStart=/usr/bin/python3 -m lauschomat.web.main --config /etc/lauschomat/web_config.yaml
Restart=always
RestartSec=5
WorkingDirectory=/var/lib/lauschomat
Environment=PYTHONUNBUFFERED=1

# Hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true

[Install]
WantedBy=multi-user.target
